# source/intro/getting_started_example/supervisor/example.conf
[program:example]
command=nodejs /opt/example/src/app.js
