# source/docs/coffee-counter/supervisor/coffee-counter.conf
[program:coffee-counter]
directory=/opt/coffee-counter/src
command=python coffee-counter.py 80
